# test/supervisord.conf
[unix_http_server]
test=a
;commented out

[program:a-program-name]
test=b

[should_not_highlight]
test=c
